# ansible/roles/flask/templates/flask_app.supervisor.ini.j2
[program:{{ project_name }}]
command={{ project_virtualenv_dir }}/bin/uwsgi_start
directory={{ project_dir }}
autostart=true
autorestart=true
